# sfsroot.base.bootstrap.d/etc/systemd/system/d2-led.service
[Unit]
Description=Run D2 LED setup

[Service]
Type=oneshot
ExecStart=/usr/bin/d2-led.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
